# ops/pg-backup.cron
# /etc/cron.d/gramjob-backup — daily PostgreSQL dump
# Format: minute hour dom month dow user command
15 3 * * * root /opt/gramjob/ops/pg-backup.sh >/dev/null 2>&1
